# data/supervisor/supervisor-apache.conf
[program:apache2]
command=/opt/nexus/scripts/start-apache.sh
numprocs=1
autostart=true
autorestart=true