# data/supervisor/supervisor-mysql.conf
[program:mysqld]
command=/opt/nexus/scripts/start-mysql.sh
numprocs=1
autostart=true
autorestart=true